# src/cron_remainder/crontab
0 15 * * *  cd /app && periodic-task-run